# crontab/stdserver-update
#!/bin/bash
./home/username/gameserver update > /dev/null 2>&1
